# mqtt_porchlight.service
[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/wemoswitch-to-mqtt/wemo-mqtt-control.py

[Install]
WantedBy=multi-user.target
